# deploy/docker-compose-app.service
# /etc/systemd/system/docker-compose-app.service
[Unit]
Description=DinnerHopping Docker Compose App (frontend + backend + mongo)
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/dinnerhopping/deploy
# Build updated images before starting containers
## NOTE: The original version used the Docker CLI plugin syntax `docker compose`.
## On this host only the standalone binary `/usr/bin/docker-compose` is installed.
## Adjusted accordingly to avoid service start failure (exit code 203 / file not found).
ExecStartPre=-/usr/bin/docker-compose -f production-docker-compose.yml build --pull
ExecStart=/usr/bin/docker-compose -f production-docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f production-docker-compose.yml down
TimeoutStartSec=0
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
